# crontab
* * * * * cd /usr/src/app && bundle exec rake wakeup_notify:send_notifications
* * * * * cd /usr/src/app && bundle exec rake routine_notify:send_notifications