# crontab
*/15 20-23,0-3 * * * cd /usr/src/app && bundle exec rake wakeup_notify:send_notifications
