# crontab
*/30 5-12 * * * cd /usr/src/app && bundle exec rake wakeup_notify:send_notifications
